# clean/crontab
23 05 * * * /var/www/html/mobilewebcam/clean/cleanup.sh
33 05 * * * /var/www/html/glider/wvsc/webcam/mobilewebcam-serverside/clean/cleanup.sh
