# clean/crontab
23 05 * * * /srv/www-html/mobilewebcam/clean/cleanup.sh

